# system/systemd/etc/systemd/user/user-applications.target
[Unit]
Description=Active when startx is called
Requires=default.target
After=default.target
